# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:go-app]
command=/app/main
autostart=true
autorestart=true
stderr_logfile=/var/log/go-app.err.log
stdout_logfile=/var/log/go-app.out.log
