# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:go-app]
command=/app/main
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
